# luci-theme-lifeline/files/etc/uci-defaults/30-luci-theme-lifeline
#!/bin/sh

. /lib/functions.sh

if ! uci show luci.theme.LifeLine | grep LifeLine >/dev/null
then 

uci set luci.theme.LifeLine=rule
uci commit 

fi

exit 0
