# lifeline-connect/files/etc/uci-defaults/81-modify-wan
#!/bin/ash

. /lib/functions.sh

if [ "uci get network.wan.metric`" != "1" ]
then 
        uci set network.wan.metric=1
        uci commit 
        /etc/init.d/network reload
fi

exit 0 

